# deploy/systemd/aleksiz-token-usage-sync.timer
[Unit]
Description=Periodically pull the redacted TokenUsage snapshot

[Timer]
OnBootSec=10m
OnUnitActiveSec=6h
RandomizedDelaySec=20m
Persistent=true
Unit=aleksiz-token-usage-sync.service

[Install]
WantedBy=timers.target
