# deploy/systemd/aleksiz-notes-sync.timer
[Unit]
Description=Periodically rebuild Aleksi Blog with fresh public Obsidian notes

[Timer]
OnBootSec=15m
OnUnitActiveSec=6h
RandomizedDelaySec=20m
Persistent=true
Unit=aleksiz-notes-sync.service

[Install]
WantedBy=timers.target
