# deploy/systemd/aleksiz-notes-sync.service
[Unit]
Description=Rebuild Aleksi Blog with fresh public Obsidian notes
After=network-online.target aleksiz-astro.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/local/sbin/aleksiz-notes-sync
TimeoutStartSec=180
UMask=0077
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/var/www/aleksiz
